# linux/cron.hourly/nmon.sh
#!/usr/bin/env bash
nmon -f -t 10 -EJ -s 60 -c 60 -m /data/nmonlog
